# anex-files/Archivos/crontab/cronfw
30 07,11,13,15,21 * * *  /root/reset-firewall.sh
